# stacks/public/messenger/service/messenger/run
#!/bin/bash

exec nodejs /var/www/app.js

sleep 5