# base-images/elasticsearch/services/elasticsearch/run
#!/bin/bash

chown -R elasticsearch:elasticsearch /usr/share/elasticsearch/data

exec setuser elasticsearch elasticsearch