# stacks/manager/proxy/service/gen/run
#!/bin/bash

set -e

exec docker-gen -watch -only-exposed -config config.ini > /dev/null 2>&1