# stacks/public/proxy/service/gen/run
#!/bin/bash

exec docker-gen     -notify "sv reload nginx" \
					-watch -only-exposed \
					/etc/service/nginx/config/nginx_conf.tmpl \
					/etc/service/nginx/config/nginx.conf