# base-images/nginx_php/services/php-fpm/run
#!/bin/bash

exec php5-fpm -c /etc/service/phpfpm/config -F -y /etc/service/phpfpm/config/php-fpm.conf